# k8s-deploy/playbooks/roles/cpu-performance/templates/lib/systemd/system/cpu-performance.service
[Unit]
Description=set cpu performance mode

[Service]
Type=forking
ExecStart=/usr/local/bin/set-cpu-performance.sh

[Install]
WantedBy=multi-user.target
